# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
# user=www-data
command=/usr/sbin/httpd -D FOREGROUND -k start